# deploy/prohub-api.service
[Unit]
Description=ProHub API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=prohub
Group=prohub
WorkingDirectory=/opt/prohub-bot
EnvironmentFile=/etc/prohub-bot.env
ExecStart=/opt/prohub-bot/venv/bin/python api/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
